# infra/files/upstart.conf
description "Monitor"
author  "Hugo <hugo@gmail.com>"
start on runlevel [2345]
stop on runlevel [016]
exec /monitor/start.sh
start on startup
