# infra/files/upstart.conf
description "Monitor"
author  "Hugo <hugo@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

script

    chdir /monitor/
    exec /monitor/start.sh

end script
